# deploy/supervisor/hablas.conf
; Background workers for hablas — one file, two programs, matching the house
; pattern (see zero.conf on the droplet).
;   hablas-queue     — database queue worker (DailyDigestNotification et al.)
;   hablas-scheduler — Laravel scheduler (SendDailyDigests daily 08:00)
;
; Install: staged to /tmp/hablas-supervisor.conf, then (root) moved into place:
;   sudo mv /tmp/hablas-supervisor.conf /etc/supervisor/conf.d/hablas.conf
;   sudo supervisorctl reread && sudo supervisorctl update
; The deploy workflow restarts these with:
;   supervisorctl restart hablas-queue:* hablas-scheduler:*
[program:hablas-queue]
process_name=%(program_name)s_%(process_num)02d
command=php /home/deploy/hablas/artisan queue:work database --sleep=3 --tries=3 --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=deploy
numprocs=1
redirect_stderr=true
stdout_logfile=/home/deploy/hablas/storage/logs/queue.log
stopwaitsecs=3600

[program:hablas-scheduler]
process_name=%(program_name)s_%(process_num)02d
command=php /home/deploy/hablas/artisan schedule:work
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=deploy
numprocs=1
redirect_stderr=true
stdout_logfile=/home/deploy/hablas/storage/logs/scheduler.log
stopwaitsecs=60
