# deploy/supervisor/hablas-queue.conf
; Database queue worker for hablas — processes DailyDigestNotification and any
; other queued jobs (QUEUE_CONNECTION=database).
; Install to /etc/supervisor/conf.d/hablas-queue.conf, then:
;   sudo supervisorctl reread && sudo supervisorctl update
; The deploy workflow restarts this group with: supervisorctl restart hablas-queue:*
;
; stopwaitsecs is high so an in-flight job finishes on a graceful --max-time
; restart rather than being killed mid-send.
[program:hablas-queue]
process_name=%(program_name)s_%(process_num)02d
command=php /home/deploy/hablas/artisan queue:work --sleep=3 --tries=3 --max-time=3600
directory=/home/deploy/hablas
autostart=true
autorestart=true
user=deploy
numprocs=1
redirect_stderr=true
stdout_logfile=/home/deploy/hablas/storage/logs/queue.log
stopwaitsecs=3600
