# deploy/supervisor/hablas-scheduler.conf
; Laravel scheduler for hablas — runs SendDailyDigests (daily 08:00).
; Install to /etc/supervisor/conf.d/hablas-scheduler.conf, then:
;   sudo supervisorctl reread && sudo supervisorctl update
; The deploy workflow restarts this group with: supervisorctl restart hablas-scheduler:*
;
; supervisord opens the logfile as root before dropping to user=deploy, so after
; the first start you may need: chown deploy:www-data storage/logs/scheduler.log
[program:hablas-scheduler]
process_name=%(program_name)s
command=php /home/deploy/hablas/artisan schedule:work
directory=/home/deploy/hablas
autostart=true
autorestart=true
user=deploy
redirect_stderr=true
stdout_logfile=/home/deploy/hablas/storage/logs/scheduler.log
stopwaitsecs=60
